# eldon/deploy/systemd/openclaw-optimizer.timer
[Unit]
Description=Run OpenClaw optimizer nightly at 02:00
Requires=openclaw-optimizer.service

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
