# eldon/gateway/infra/openclaw-gateway.service
[Unit]
Description=OpenClaw Lola WhatsApp Gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/openclaw/eldon/gateway

EnvironmentFile=/opt/openclaw/.env

ExecStart=/opt/openclaw/.venv/bin/python scripts/run_gateway.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openclaw-gateway

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/openclaw/.lola
ReadWritePaths=/opt/openclaw/eldon/logs

[Install]
WantedBy=multi-user.target
